# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx.log
stderr_logfile=/var/log/nginx.log

[program:vllm-chat]
command=/opt/vllm-env/bin/python -m vllm.entrypoints.openai.api_server
    --model              /models/qwen2.5-coder-32b-awq
    --host               127.0.0.1
    --port               8000
    --served-model-name  qwen2.5-coder-32b
    --quantization       awq_marlin
    --max-model-len      32768
    --gpu-memory-utilization 0.75
    --enable-prefix-caching
autostart=true
autorestart=true
stdout_logfile=/var/log/vllm-chat.log
stdout_logfile_maxbytes=50MB
stderr_logfile=/var/log/vllm-chat.log
stderr_logfile_maxbytes=0

[program:vllm-complete]
command=/bin/bash /opt/start_3b.sh
autostart=true
autorestart=true
stdout_logfile=/var/log/vllm-complete.log
stdout_logfile_maxbytes=50MB
stderr_logfile=/var/log/vllm-complete.log
stderr_logfile_maxbytes=0
